# rootfs/etc/init.d/spark-server
#!/sbin/openrc-run

mkdir -p /var/run
mkdir -p /var/log
SCRIPT="/usr/bin/spark-server"
PIDFILE=/var/run/spark-server.pid
LOGFILE=/var/log/spark-server.log

start() {
    if [ -f $PIDFILE ] && [ -s $PIDFILE ] && kill -0 $(cat $PIDFILE); then
        echo "Server is already running" >&2
        return 1
    fi

    echo "Starting the Spark server process" >&2
    local CMD="$SCRIPT &> \"$LOGFILE\" & echo \$!"
    local CMD="$SCRIPT &> \"$LOGFILE\" & echo \$!"

    su -c "$CMD" root > "$PIDFILE"

    sleep 2

    PID=$(cat $PIDFILE)
    if pgrep -f "spark-server" > /dev/null
    then
        echo "spark-server is now running under PID $PID"
    else
        echo "Failed to start the spark server"
    fi


    return 0
}

stop() {
    if [ ! -f "$PIDFILE" ] || ! kill -0 $(cat "$PIDFILE"); then
        echo 'Service not running' >&2
        return 1
    fi
    echo 'Stopping service…' >&2
    kill -15 $(cat "$PIDFILE") && rm -f "$PIDFILE"
    echo 'Service stopped' >&2
}